# mvp1/air-mouse-leds.service
[Unit]
Description=Air Mouse LED Control Daemon
Documentation=https://github.com/sant0s12/openrgb-air-mouse
After=graphical-session.service
Wants=graphical-session.service

[Service]
Type=simple
User=sant
Environment="DISPLAY=:0"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus"
ExecStart=/usr/bin/python3 /home/sant/Área de trabalho/PROJETOS/openrgb/controle_led.py
Restart=on-failure
RestartSec=5

# Log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=air-mouse-leds

[Install]
WantedBy=graphical.target
